# sparse/etc/profile.d/hide_stowaways_from_find.sh
find() {
    p=${1:-.}
    /bin/find "$p" ! \( -path "${p/%\//}/data/.stowaways" -prune \) "${@:2}"
}
